# roles/kong/templates/kong.service.j2
[Unit]
Description=Run script at startup after network becomes reachable

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/root/kong_start.sh
TimeoutStartSec=0

[Install]
WantedBy=default.target
